# WebApiServer/PluginLoader.service
[Unit]
Description=Plugin Loader server

[Service]
Type=notify
ExecStart=

[Install]
WantedBy=multi-user.target